# targets/nanopi_overlay/etc/udev/rules.d/10-brcm-bluetooth.rules
SUBSYSTEM=="rfkill", ACTION=="change", ATTR{name}=="bcm20710 Bluetooth", ATTR{state}=="1", RUN+="/bin/systemctl --no-block start brcm-bluetooth.service"
SUBSYSTEM=="rfkill", ACTION=="change", ATTR{name}=="bcm20710 Bluetooth", ATTR{state}=="0", RUN+="/bin/systemctl --no-block stop brcm-bluetooth.service"
