# targets/nanopi_overlay/etc/systemd/system/brcm-bluetooth.service
[Unit]
Description=Broadcom firmware and hci service

[Service]
Type=simple
ExecStart=/usr/local/bin/brcm-bluetooth.sh
