# targets/mini210_overlay/etc/profile.d/ts_env.sh
export TERM=xterm
export TSLIB_CALIBFILE=/etc/pointercal
export TSLIB_CONFFILE=/etc/ts.conf
export TSLIB_TSDEVICE=/dev/touchscreen-1wire
export TSLIB_FBDEVICE=/dev/fb0
export TSLIB_CONSOLEDEVICE=/dev/tty1
export TSLIB_PLUGINDIR=/usr/lib/ts
